# etc/supervisord.d/kafka.ini
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
